# scripts/crontab
# scripts/crontab
# Backup ogni notte alle 02:00
0 2 * * * DB_USER=$DB_USER DB_NAME=$DB_NAME /backup.sh >> /var/log/backup.log 2>&1
